# config/di_3dp_udp_renderer.conf
[program:di3dp]
command=python -m disinfo.renderers.udp --fps 24
directory=/home/pi/disinfo
user=pi
environment=DI_CONFIG_PATH=".config-3dpanel.json"
