# config/di_dataservice.conf
[program:didataservice]
command=/home/noop/run/disinfo/.venv/bin/python -m disinfo.drat.data_service
directory=/home/noop/run/disinfo
user=noop
