# config/di_dataservice.conf
[program:didataservice]
command=/bin/bash -c "source .env && python data_service.py"
directory=/home/pi/disinfo
user=pi
