# config/di_dataservice.conf
[program:didataservice]
command=python -m disp_info.data_service
directory=/home/pi/disinfo
user=pi
