# config/di_dataservice.conf
[program:didataservice]
command=python -m disinfo.drat.data_service
directory=/home/pi/disinfo
user=pi
