# config/di_dataservice.conf
[program:didataservice]
command=watchmedo auto-restart -d disinfo/ --pattern="*.py" --recursive -- python -m disinfo.drat.data_service
directory=/home/pi/disinfo
user=pi
