# config/di_dataservice.conf
[program:didataservice]
command=watchmedo auto-restart -d disp_info/ --pattern="*.py" --recursive -- python -m disp_info.drat.data_service
directory=/home/pi/disinfo
user=pi
