# config/di_haservice.conf
[program:dihaservice]
command=watchmedo auto-restart -d disinfo/ --pattern="*.py" --recursive -- python -m disinfo.drat.ha_service
directory=/home/pi/disinfo
user=pi
