# config/di_haservice.conf
[program:dihaservice]
command=python -m disinfo.drat.ha_service
directory=/home/pi/disinfo
user=pi
