# config/di_haservice.conf
[program:dihaservice]
command=python -m disp_info.ha_service
directory=/home/pi/disinfo
user=pi
