# config/di_haservice.conf
[program:dihaservice]
command=watchmedo auto-restart -d disp_info/ --pattern="*.py" --recursive -- python -m disp_info.drat.ha_service
directory=/home/pi/disinfo
user=pi
