# config/di_haservice.conf
[program:dihaservice]
command=/home/noop/.pyenv/versions/dipy/bin/python -m disinfo.drat.ha_service
directory=/home/noop/run/disinfo
user=noop
