# config/di_shazam.conf
[program:dishazamservice]
command=/home/noop/.pyenv/versions/dipy/bin/python -m disinfo.drat.shazam
directory=/home/noop/run/disinfo
user=noop
