# config/di_salon_udp_renderer.conf
[program:disalon]
command=python -m disinfo.renderers.udp --fps 33
directory=/home/pi/disinfo
user=pi
environment=DI_CONFIG_PATH=".config-salon.json"
